# .config/systemd/user/timers.target.wants/mbsync.timer
[Unit]
Description=Run mbsync every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
