# .config/systemd/user/default.target.wants/keyledsd-dunst-forwarder.service
[Unit]
Description=Forward dunst notifications to keyledsd via DBus

[Service]
ExecStart=%h/.local/bin/keyledsd-dunst-forwarder
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
